# systemd/claude-memory-manager.service
[Unit]
Description=Claude Memory Manager Daemon (Rust)
Documentation=https://github.com/recodeee/claude-memory-manager-daemon
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/Documents/claude-memory-manager-daemon
ExecStart=%h/Documents/claude-memory-manager-daemon/target/release/cmmd run
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=30
TimeoutStopSec=15

# Default-safe behavior; flip DRY_RUN=false in an override only after you trust it.
Environment=DRY_RUN=true
Environment=CMMD_LOG=info

# Make sure the daemon can find `claude` and `authmux`. Systemd user units
# inherit a minimal PATH that does NOT include ~/.nvm/.../bin, so we extend it.
Environment=PATH=%h/.nvm/versions/node/v22.22.0/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAUDE_BIN=%h/.nvm/versions/node/v22.22.0/bin/claude
Environment=AUTHMUX_BIN=%h/.nvm/versions/node/v22.22.0/bin/authmux

StandardOutput=append:%t/claude-memory-manager.log
StandardError=append:%t/claude-memory-manager.log

# Hardening — daemon needs read on $HOME, write on memory dir + /tmp.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.claude/projects/-home-deadpool/memory /tmp
PrivateTmp=false
NoNewPrivileges=true

[Install]
WantedBy=default.target
